# data/setupTemplates/pz-backup.service
[Unit]
Description=Project Zomboid - Hourly Data Backup
After=zomboid.service

[Service]
Type=oneshot
ExecStart=/home/pzuser/pzmanager/scripts/backup/dataBackup.sh
StandardOutput=journal
StandardError=journal
