# data/setupTemplates/pz-maintenance.timer
[Unit]
Description=Project Zomboid - Daily Maintenance Timer

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
