# data/setupTemplates/pz-restart.service
[Unit]
Description=Project Zomboid - Redémarrage planifié (stabilité mémoire)

[Service]
Type=oneshot
# `restart 15m` bloque ~15 min (avertissements joueurs en jeu via sleep) :
# sans ce timeout large, systemd tuerait le service à DefaultTimeoutStartSec (90s).
TimeoutStartSec=1500
ExecStart=%h/pzmanager/pzm server restart 15m --automatic --reason "Redémarrage planifié (stabilité)"
StandardOutput=journal
StandardError=journal
